# myconfig/autostart/conky.desktop
[Desktop Entry]
Name=Conkyrc
Comment=conky
Icon=application-x-executable
Exec=conky -q
Terminal=false
Type=Application
