# myconfig/autostart/redshift.desktop
[Desktop Entry]
Name=redshift
Comment=redshift
Icon=application-x-executable
Exec=redshift -l 50:15 -t 5500:4200
Terminal=false
Type=Application
