# myconfig/autostart/thunar.desktop
[Desktop Entry]
Name=Thunar
Comment=Thunar deamon
Icon=Thunar
Exec=thunar --daemon
Terminal=false
Type=Application
