# myconfig/autostart/volumed.desktop
[Desktop Entry]
Name=Volumed Xfce4
Comment=Volumed
Icon=application-x-executable
Exec=xfce4-volumed
Terminal=false
Type=Application
